# standalone-ng/s6-overlay/s6-rc.d/obwatcher-log/run
#!/command/execlineb -P
s6-log -b -- /var/log/jam-ng/obwatcher
