# standalone/etc/cron.d/jam-savelog
# rotate every second day
0 0 */2 * * root /usr/local/bin/jam-savelog.sh
