# standalone/etc/cron.d/jam-savelog
# daily
0 0 * * * root /usr/local/bin/jam-savelog.sh
